# cron.d/batch.cron
0 9,17 * * * cd /bootroomchat/XTScore/XTScore; /usr/local/bin/python3.8 -m sync_oss save && /usr/local/bin/python3.8 -m opta_match_map && /usr/local/bin/python3.8 -m xt_run calculate_and_generate_report > /bootroomchat/mylog.log 2>&1 &
